# nginx/supervisor/whalefall-dev.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/

[program:whalefall]
command=/app/.venv/bin/gunicorn --config /app/gunicorn.conf.py app:app
directory=/app
user=root
autostart=true
autorestart=true
redirect_stderr=false
stdout_logfile=/app/userdata/logs/whalefall.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/app/userdata/logs/whalefall_error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10

# 开发环境不需要supervisor管理Nginx，使用系统服务
